# tools/com.reaper.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.reaper.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rus/Projects/ReaperOS/.reaper_venv/bin/python3</string>
        <string>-u</string>
        <string>/Users/rus/Projects/ReaperOS/tools/watchdog.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rus/Projects/ReaperOS</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/rus/Projects/ReaperOS/logs/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rus/Projects/ReaperOS/logs/watchdog.log</string>
</dict>
</plist>
